# checkoff-printer.service
[Unit]
Description=Checkoff Printer Web Service
After=network.target

[Service]
User=root
WorkingDirectory=/home/printer/checkoff-printer
Environment="PRINTER_MODE=usb"
# Use uv to run the application from the WorkingDirectory
ExecStart=/home/printer/.local/bin/uv run app.py
Restart=always

[Install]
WantedBy=multi-user.target
