# checkoff-printer.service
[Unit]
Description=Checkoff Printer Web Service
After=network.target

[Service]
User=printer
WorkingDirectory=/home/printer/checkoff-printer
Environment="PRINTER_MODE=usb"
# Use uv to run the application
ExecStart=/home/printer/.local/bin/uv run --project /home/printer/checkoff-printer /home/printer/checkoff-printer/app.py
Restart=always

[Install]
WantedBy=multi-user.target
