# checkoff-printer.service
[Unit]
Description=Checkoff Printer Web Service
After=network.target

[Service]
User=printer
WorkingDirectory=/home/printer/checkoff-printer
Environment="PRINTER_MODE=usb"
# Adjust path to venv python (check if it's .venv or venv)
ExecStart=/home/printer/checkoff-printer/venv/bin/python /home/printer/checkoff-printer/app.py
Restart=always

[Install]
WantedBy=multi-user.target
